# ipv6/tayga/files/tayga.init
#!/bin/sh /etc/rc.common

START=80
STOP=15

start() {
	tayga --mktun
	ip link set nat64 up
	ip addr add 2001:470:1f15:982::2 dev nat64  
	ip addr add 192.168.1.212 dev nat64
	ip route add 64:ff9b::/96 dev nat64  
	ip route add 192.168.1.16/30 dev nat64      
	echo 1 > /proc/sys/net/ipv4/conf/eth0/proxy_arp
	tayga
}

stop() {
	killall tayga
}

