# net/osnmpd/files/osnmpd.init
#!/bin/sh /etc/rc.common

START=80
STOP=40

start() {
	/usr/sbin/snmpd -c /etc/snmp/snmpd.conf
}

stop() {
	killall snmpd
}

